# debian/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

LOCALSCANCONF=/etc/exim4/conf.d/main/15_sa-exim_plugin_path
DISABLELOCALSCANTEXT='# sa-exim plugin is disabled, because it was uninstalled'
DISABLELOCALSCANMD5=`echo "${DISABLELOCALSCANTEXT}" | md5sum | cut -d\  -f1`

OLDCRONJOB=/etc/cron.hourly/greylistclean
OLDCRONJOBREN=/etc/cron.hourly/.greylistclean.insecure.deleteme
if test -f $OLDCRONJOB; then
    echo "disabling old insecure $OLDCRONJOB"
    echo "(renamed to $OLDCRONJOBREN)"
    echo "See new cronjob in /etc/cron.d/greylistclean"
    mv $OLDCRONJOB $OLDCRONJOBREN
fi

case "$1" in
    configure)
	if [ ! -e /var/spool/sa-exim ] ; then
		# Debian-exim should exist as we depend on exim4-base
		install -d -m770 -oDebian-exim -gDebian-exim \
		/var/spool/sa-exim
	fi
	# Support for greylisting tuplets (written by default as Debian-exim by spamd)
	if [ ! -e /var/spool/sa-exim/tuplets ] ; then
		install -d -m770 -oDebian-exim -gDebian-exim \
		/var/spool/sa-exim/tuplets
	elif [ `stat -c "%U:%G" /var/spool/sa-exim/tuplets` = "nobody:Debian-exim" ] &&
		dpkg --compare-versions "$2" lt-nl "4.2.1-14"; then
	    	# If ownership is as before 4.2.1-14, give the Debian-exim
	    	# group write permissions. Bug 563492.
		chmod g+w /var/spool/sa-exim/tuplets
	fi
	# clean up temporary file generated by postrm uninstall
	if [ -e "${LOCALSCANCONF}.rul" ] && \
	   [ "`md5sum ${LOCALSCANCONF}.rul | cut -d\  -f1`" = "${DISABLELOCALSCANMD5}" ]; then
	   rm ${LOCALSCANCONF}.rul
	fi
	if [ -x /etc/init.d/exim4 ]; then
		if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
			invoke-rc.d exim4 reload || true
		else
			/etc/init.d/exim4 reload || true
		fi
	fi
	;;

esac

#DEBHELPER#
